# studyguide.service
[Unit]
Description=StudyGuide Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/studyguide
# Use the venv's python so all dependencies are found
ExecStart=/home/ubuntu/studyguide/venv/bin/python /home/ubuntu/studyguide/bot.py
EnvironmentFile=/home/ubuntu/studyguide/.env
Restart=always
RestartSec=5
StandardOutput=append:/home/ubuntu/studyguide/bot.log
StandardError=append:/home/ubuntu/studyguide/bot.log

[Install]
WantedBy=multi-user.target